# docker-udev-aliases-sync-script/templates/99-docker-test.rules.j2
{% for device in devices %}
# Map {{ device.kernel_address }} to /dev/{{ device.symlink }}
ACTION=="add|change", SUBSYSTEM=="tty", SUBSYSTEMS=="usb", KERNELS=="{{ device.kernel_address }}", RUN+="/usr/local/bin/docker-static-usb-nodes.sh add /dev/{{ device.symlink }} %M %m"
ACTION=="remove", SUBSYSTEM=="tty", SUBSYSTEMS=="usb", KERNELS=="{{ device.kernel_address }}", RUN+="/usr/local/bin/docker-static-usb-nodes.sh remove /dev/{{ device.symlink }} %M %m"

{% endfor %}
